# deb/usr/lib/systemd/user/allrelay.service
[Unit]
Description=AllRelay - Android Phone Media Stream Server
After=pipewire.service pipewire-pulse.service

[Service]
Type=simple
ExecStart=/usr/bin/allrelay-server --web-port 9090 -v
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
