# deb/DEBIAN/postinst
#!/bin/bash
set -e

# Load v4l2loopback module if not already loaded
if ! lsmod | grep -q v4l2loopback; then
    modprobe v4l2loopback video_nr=10 card_label="AllRelay Cam" exclusive_caps=1 || true
fi

# Ensure v4l2loopback auto-loads on boot
echo "v4l2loopback" > /etc/modules-load.d/allrelay.conf
echo "options v4l2loopback video_nr=10 card_label=\"AllRelay Cam\" exclusive_caps=1" > /etc/modprobe.d/allrelay.conf

# Remove old system-wide service if present (migrated to user service)
if [ -f /etc/systemd/system/allrelay.service ]; then
    systemctl stop allrelay 2>/dev/null || true
    systemctl disable allrelay 2>/dev/null || true
    rm -f /etc/systemd/system/allrelay.service
    systemctl daemon-reload || true
fi

# Both services are on-demand. App Grid (`allrelay tray`) starts the backend
# and tray together; neither service should start automatically at login.
# Disable old global enablement when installing or upgrading.
systemctl --global disable allrelay.service allrelay-tray.service || true

# Do not leave an already-running copy behind after an upgrade.
if [ -n "${SUDO_USER:-}" ] && [ "$SUDO_USER" != "root" ]; then
    su - "$SUDO_USER" -c "XDG_RUNTIME_DIR=/run/user/\$(id -u) systemctl --user daemon-reload" || true
    su - "$SUDO_USER" -c "XDG_RUNTIME_DIR=/run/user/\$(id -u) systemctl --user stop allrelay allrelay-tray" || true
fi

echo ""
echo "=============================================="
echo "  AllRelay installed!"
echo "  Tray controls:  allrelay tray"
echo "  Open dashboard: allrelay open"
echo "  Show URL:       allrelay url"
echo "  Virtual camera: /dev/video10 (AllRelay Cam)"
echo ""
echo "  Start/stop: systemctl --user [start|stop|restart] allrelay"
echo "  Logs:       journalctl --user -u allrelay -f"
echo "=============================================="
