# deb/DEBIAN/prerm
#!/bin/bash
set -e

systemctl stop allrelay || true
systemctl disable allrelay || true
systemctl daemon-reload || true
