# deb/etc/systemd/system/allrelay.service
[Unit]
Description=AllRelay - Android Phone Media Stream Server
After=network.target pulseaudio.service

[Service]
Type=simple
User=%u
Environment=HOME=%h
Environment=XDG_RUNTIME_DIR=/run/user/%U
ExecStart=/usr/bin/allrelay-server --web-port 9090 -v
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
