# configs/systemd/allrelay.service
# AllRelay Server — systemd service
#
# Runs the AllRelay Ubuntu server that connects to the Android phone
# and routes all streams. Requires the phone IP to be configured.
#
# Install:
#   sudo cp allrelay.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now allrelay
#
# The service runs per-user to access PipeWire (user session).

[Unit]
Description=AllRelay Server — Wireless phone peripherals for Ubuntu
After=network-online.target pipewire.service pipewire-pulse.service
Wants=network-online.target
Requires=pipewire.service

[Service]
Type=simple
# v4l2loopback setup (for camera stream, Phase 3)
ExecStartPre=/sbin/modprobe v4l2loopback devices=1 video_nr=10 card_label="AllRelay Camera" exclusive_caps=1
# Start the AllRelay server
# PHONE_IP must be set in /etc/allrelay/phone_ip or via environment
EnvironmentFile=-/etc/allrelay/phone_ip
ExecStart=/usr/local/bin/allrelay-server --host ${PHONE_IP:-192.168.1.100}
ExecStopPost=/sbin/rmmod v4l2loopback 2>/dev/null || true
Restart=always
RestartSec=5

# Run as the local user to access PipeWire
User=%u
Group=%u

# PipeWire environment
Environment=PIPEWIRE_RUNTIME_DIR=/run/user/%U
Environment=XDG_RUNTIME_DIR=/run/user/%U

[Install]
WantedBy=multi-user.target
